# udev/70-asteroid-docking-bay.rules
# SPDX-License-Identifier: GPL-3.0-only
# SPDX-FileCopyrightText: 2026 Timo Könnecke (moWerk) <mo@mowerk.net>
# SPDX-FileCopyrightText: 2023 Ed Beroset <beroset@ieee.org>
# udev rules for asteroid-docking-bay
# ------------------------------------
# Grants non-root access to smart USB hubs (for uhubctl) and AsteroidOS
# watches (for adb).
#
# Access is granted two ways, and each rule below uses both so it works on
# desktop and headless/SSH setups alike:
#   TAG+="uaccess"     hands the device to the user of the active local
#                      login seat (works on a normal desktop login, and on
#                      a machine with console autologin).
#   GROUP="users"      gives every member of the 'users' group plain rw
#                      access via file permissions. This is what makes it
#                      work over SSH, where there is no active local seat.
#
# NOTE: earlier versions used GROUP="plugdev". Many distros (Arch,
# EndeavourOS, Fedora) do not create a 'plugdev' group, in which case udev
# silently cannot assign it and the device stays root-only. 'users' exists
# almost everywhere; if your distro does not put you in it, pick a group you
# are in (check with `groups`) or add yourself to 'users'.
#
# Usage:
#   sudo cp 70-asteroid-docking-bay.rules /etc/udev/rules.d/
#   sudo udevadm control --reload-rules && sudo udevadm trigger --action=add
#   # verify (no sudo): uhubctl -l    should now list your hub
#
# If a hub that is already plugged in does not pick up the new access after
# the trigger, unplug and replug its upstream cable once.
#
# Finding your hub's Vendor/Product ID:
#   lsusb | grep -i hub
#   uhubctl -l            # shows hub locations and IDs
#
# Add one SUBSYSTEM line per smart hub. Confirmed and common examples below.

# ── Smart USB hubs (for uhubctl) ──────────────────────────────────────────────

# Vendor ID 17ef = Lenovo (ThinkPad USB 3.0 Ultra Dock DK-1523, confirmed true VBUS switching)
SUBSYSTEM=="usb", ATTR{idVendor}=="17ef", MODE="0664", GROUP="users", TAG+="uaccess"

# Vendor ID 0bda = Realtek (RTS5411, e.g. 16-port RSHTech-style hubs; confirmed
# per-port power switching in the hub descriptor. Note: the RTS5411 reports
# "Generic USB2.1 Hub", so uhubctl needs this rule to reach it by vendor ID.)
SUBSYSTEM=="usb", ATTR{idVendor}=="0bda", MODE="0664", GROUP="users", TAG+="uaccess"

# Vendor ID 2109 = VIA Labs (common in USB 3 hubs)
# SUBSYSTEM=="usb", ATTR{idVendor}=="2109", MODE="0664", GROUP="users", TAG+="uaccess"

# Vendor ID 0424 = Microchip/SMSC (common in Pi hats and dev-board hubs)
# SUBSYSTEM=="usb", ATTR{idVendor}=="0424", MODE="0664", GROUP="users", TAG+="uaccess"

# Vendor ID 05e3 = Genesys Logic / ALCOR (data-line disconnect only, not true VBUS)
# SUBSYSTEM=="usb", ATTR{idVendor}=="05e3", MODE="0664", GROUP="users", TAG+="uaccess"

# ── Rootless port power switching (direct sysfs write) ────────────────────────
# asteroid-docking-bay switches port power by writing the kernel's per-port
# power attribute directly — near-instant, with no libusb bus scan (uhubctl's
# per-command scan is what made switching slow and churned the bus):
#     /sys/bus/usb/devices/<hub>:1.0/<hub>-portN/disable   (0 = on, 1 = off)
# That attribute is root-only by default; without write access the tool falls
# back to slow uhubctl. These rules grant the 'users' group write on the port
# attributes of the smart hubs above.
#
# A USB *port* is not itself a udev device (no uevent), so we can't match it
# directly — instead we match each hub on 'add' and chmod its child port
# attributes. Match the SAME vendors you enabled for uhubctl above; the host's
# own root hub is a different vendor and is never touched. Add a line per
# vendor you use.
#
# After installing these rules run (or just replug the hub's upstream cable):
#     sudo udevadm control --reload-rules
#     sudo udevadm trigger --action=add
# so the RUN fires with the ports already enumerated. Verify (no sudo):
#     ls -l /sys/bus/usb/devices/*:1.0/*-port*/disable   # group 'users', -rw-rw-r--
#
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0bda", RUN+="/bin/sh -c 'for f in /sys%p/*:1.0/*-port*/disable; do test -e $f && chgrp users $f && chmod g+w $f; done'"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="17ef", RUN+="/bin/sh -c 'for f in /sys%p/*:1.0/*-port*/disable; do test -e $f && chgrp users $f && chmod g+w $f; done'"
# ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="2109", RUN+="/bin/sh -c 'for f in /sys%p/*:1.0/*-port*/disable; do test -e $f && chgrp users $f && chmod g+w $f; done'"

# ── AsteroidOS watches (for adb) ──────────────────────────────────────────────
# AsteroidOS uses the Android gadget framework, so standard Android rules apply.
# On Arch Linux: install 'android-udev' from AUR.
# On Debian/Ubuntu: install 'android-sdk-platform-tools' or add rules manually.
#
# Generic fallback (covers most watches; prefer the android-udev package):
# SUBSYSTEM=="usb", ATTR{idVendor}=="18d1", MODE="0664", GROUP="users", TAG+="uaccess"
